# agent/snapd.service
[Unit]
Description=Snapshot Daemon - Blockchain Node Monitoring and Backup Service
Documentation=https://github.com/nodexeus/chainsnaps
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=snapd
Group=snapd

# Path to the snapd binary
ExecStart=/usr/local/bin/snapd --config /etc/snapd/config.yaml

# Restart policy
Restart=on-failure
RestartSec=10s

# Environment variables
# Database password should be set here or in an environment file
# Environment="DB_PASSWORD=your_secure_password"
# Or use an environment file:
EnvironmentFile=-/etc/snapd/environment

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/snapd

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=snapd

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
